# system_files/system_files/usr/lib/systemd/system/dconf-update.service
[Unit]
Documentation=https://github.com/coreos/rpm-ostree/issues/1944

[Service]
Type=oneshot

[Install]
WantedBy=multi-user.target
